# rootfs/etc/s6/game/run
#!/bin/bash

if [ -f ./setup ]
then
  source ./setup
fi

STARTCMD="/usr/bin/launcher"
exec ${STARTCMD}
